# deploy/oxy-api.service
[Unit]
Description=OXY API backend
After=network.target

[Service]
Type=simple
User=oxy
Group=oxy
WorkingDirectory=/opt/oxy/api
EnvironmentFile=/opt/oxy/api/.env
ExecStart=/usr/bin/node /opt/oxy/api/dist/index.js
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
